# root/home/pi/.config/autostart/tint2.desktop
[Desktop Entry]
Name=Autostart Tint2 panel on boot
Exec=/bin/bash -c 'sleep 5 && tint2'
