# root/home/pi/.config/autostart/windows-10-theme.desktop
[Desktop Entry]
Name=Autostart Tint2 panel on boot
Exec=/bin/bash -c '~/Windows-10/runonstartup'
